# packaging/systemd/ttygated.service
[Unit]
Description=ttygate terminal gateway
Documentation=https://github.com/tkolsto/ttygate
After=network.target
StartLimitIntervalSec=30s
StartLimitBurst=5

[Service]
Type=notify
NotifyAccess=main
WatchdogSec=6s
ExecStart=/usr/local/bin/ttygated /etc/ttygate/ttygate.toml
User=ttygate
Group=ttygate
UMask=0077

StateDirectory=ttygate
StateDirectoryMode=0700
LogsDirectory=ttygate
LogsDirectoryMode=0700

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateDevices=yes
PrivateTmp=yes
PrivateMounts=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
DevicePolicy=closed

KillMode=control-group
TimeoutStopSec=15s
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
